# .config/systemd/user/niri.service.d/ssh-agent.conf
[Unit]
Wants=gcr-ssh-agent.socket
After=gcr-ssh-agent.socket
